# packaging/lib/systemd/system/goddns.service
[Unit]
Description=goddns

[Service]
User=goddns
ExecStart=/usr/local/bin/goddns
Type=oneshot
