# packaging/lib/systemd/system/goddns.timer
[Unit]
Description=goddns

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target